# upstart.conf
description "Startup Script for node google radio"
author      "Adam Duncan"

start on (local-filesystems and net-device-up IFACE=eth0)
stop  on shutdown

respawn                # restart when job dies
respawn limit 5 60     # give up restart after 5 respawns in 60 seconds

env EMAIL=xxxx
env PASSWORD=xxxx
env AIRPROXY_HOST="192.168.0.10"
env AIRPROXY_PORT=5001

script
        exec /usr/local/bin/node /var/application/google-radio/index.js
end script
